# etc/supervisor/conf.d/sshd.conf
[program:sshd]
autorestart = true
command=/usr/sbin/sshd -E /var/log/sshd.log
