# etc/supervisor/conf.d/sssd.conf
[program:sssd]
autorestart = true
command=/usr/sbin/sssd -D -f
